# conf/crontab
SHELL=/bin/bash
*/5 * * * * cd /srv/www/techism-staging && source venv/bin/activate && export DJANGO_SETTINGS_MODULE="techism.settings.staging" && ./manage.py sendtwitter

